# howto/systemd/ba.service
[Unit]
Description=baba chess trainer

[Service]
ExecStart=/baba-chess-trainer/howto/systemd/ba.sh

[Install]
WantedBy=multi-user.target

