# howto/cloud/linode/ba.service
[Unit]
Description=baba chess trainer

[Service]
ExecStart=/baba-chess-trainer/ba.sh

[Install]
WantedBy=multi-user.target

